# services/update-feeds/run
#!/bin/bash
/sbin/setuser abc php /config/www/tt-rss/update_daemon2.php > /dev/null 2>&1
