# root/etc/services.d/update-feeds/run
#!/usr/bin/with-contenv bash

	s6-setuidgid abc php7 \
	/var/www/html/update_daemon2.php > /dev/null 2>&1
